# packages/base-conf/data.overlay.mount
[Unit]
Description=cOS tmpfs mount for the ephemeral /data

[Mount]
Where=/data.overlay
What=tmpfs
Type=tmpfs
Options=defaults,size=128M
