# packages/base-conf/etc.mount
[Unit]
Description=cOS overlayfs mount for /etc
RequiresMountsFor=/run/overlay
Requires=overlay-setup.service

[Mount]
Where=/etc
What=overlay
Type=overlay
LazyUnmount=true
Options=defaults,lowerdir=/etc,upperdir=/run/overlay/etc/upper,workdir=/run/overlay/etc/work
