# packages/base-conf/overlay-setup.service
[Unit]
Description=Create cOS overlayfs structure for rw ephemeral folders
RequiresMountsFor=/run/overlay
DefaultDependencies=no
Before=etc.mount
Before=var.mount
Before=root.mount

[Service]
Type=oneshot
ExecStart=/usr/bin/mkdir -p /run/overlay/etc/upper /run/overlay/etc/work /run/overlay/var/upper /run/overlay/var/work /run/overlay/root/upper /run/overlay/root/work

[Install]
WantedBy=local-fs.target
