# packages/base-conf/root.mount
[Unit]
Description=cOS overlayfs mount for /root
RequiresMountsFor=/run/overlay
Requires=overlay-setup.service

[Mount]
Where=/root
What=overlay
Type=overlay
LazyUnmount=true
Options=defaults,lowerdir=/root,upperdir=/run/overlay/root/upper,workdir=/run/overlay/root/work
