# packages/base-conf/var.mount
[Unit]
Description=cOS overlayfs mount for /var
RequiresMountsFor=/run/overlay
Requires=overlay-setup.service

[Mount]
Where=/var
What=overlay
Type=overlay
LazyUnmount=true
Options=defaults,lowerdir=/var,upperdir=/run/overlay/var/upper,workdir=/run/overlay/var/work
