# packages/base-conf/var.mount
[Unit]
Description=cOS overlayfs mount for /var
RequiresMountsFor=/data.overlay
Requires=overlay-setup.service

[Mount]
Where=/var
What=overlay
Type=overlay
LazyUnmount=true
Options=defaults,lowerdir=/var,upperdir=/data.overlay/var/upper,workdir=/data.overlay/var/work
