# packages/cos-setup/cos-setup-fs-after.service
[Unit]
Description=cOS system after FS setup
DefaultDependencies=no
After=local-fs.target
Wants=local-fs.target
Before=sysinit.target
RequiresMountsFor=/etc

OnFailure=emergency.target
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/cos-setup fs.after

[Install]
WantedBy=sysinit.target
